# dist/supervisor_dummy_service.conf
[program:dummy]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/signal.php
user=www-data
numprocs=4
autostart=true
autorestart=true
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/dummy-worker.log