# dist/supervisor_dummy_service_non_native.conf
[program:dummy]
process_name=%(program_name)s_%(process_num)02d
command=docker run --rm --name supervisor_%(program_name)s_%(process_num)02d -u www-data -v "/home/vagrant/code/www:/var/www/html" etable/php7.2 /usr/local/bin/php /var/www/html/signal.php
user=vagrant
numprocs=4
autostart=true
autorestart=true
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/dummy-worker.log

[program:dummy-dev]
process_name=%(program_name)s_%(process_num)02d
command=docker run --rm --name supervisor_%(program_name)s_%(process_num)02d -u www-data -v "/home/vagrant/code/www:/var/www/html" etable/php7.2-dev /usr/local/bin/php /var/www/html/signal.php
user=vagrant
numprocs=4
autostart=true
autorestart=true
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/dummy-worker.log